# packer/scripts/hhfab-vlab.service
[Unit]
Description=Hedgehog VLAB Service
Documentation=https://docs.hedgehog.cloud/latest/vlab/overview/
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service
ConditionPathExists=!/var/lib/hedgehog-lab/vlab-initialized

[Service]
Type=oneshot
ExecStart=/usr/local/bin/hhfab-vlab-runner
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
Restart=no

# Run as hhlab user for proper permissions
User=hhlab
Group=hhlab

# Environment
Environment="VLAB_WORK_DIR=/opt/hedgehog/vlab"
Environment="VLAB_TIMEOUT=1800"
Environment="LOG_FILE=/var/log/hedgehog-lab/modules/vlab.log"
Environment="TERM=xterm-256color"

# Working directory
WorkingDirectory=/opt/hedgehog/vlab

# Security hardening
# NOTE: hhfab needs to run sudo helpers to create TAP devices. Disabling
# NoNewPrivileges fixes Bug #21 (systemd prevented sudo escalation).
# PrivateTmp stays disabled so tmux sockets remain accessible to students.

[Install]
WantedBy=multi-user.target
